# systemd/restic-check@.timer
[Unit]
Description=Monthly check of restic backup

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
